# vmware-vmmon.service
[Unit]
Description=VMware Manager

[Service]
Type=forking
ExecStart=-/sbin/modprobe vmmon

[Install]
WantedBy=multi-user.target
